# modules/vulnerabilities/unix/ftp/pachev_ftp_server_1_path_traversal/templates/pachevftp.service.erb
[Unit]
Description=Pachev FTP Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=<%= @user %>
ExecStart=<%= @release_dir %>/ftp_server -p 2121
WorkingDirectory=<%= @release_dir %>/

[Install]
WantedBy=multi-user.target
